# tide_client/tidegauge.service
[Unit]
Description=Tide Gauge Service
# Keep collecting locally during offline boots; sync can catch up later.
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/tide

ExecStart=/home/pi/tide/tide_client

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
